# pkg/production/rell.service
[Unit]
Description=Rell Server
Wants=redis-socket.service

[Service]
User=rell
ExecStart=/usr/bin/rell -c /etc/conf.d/rell
Restart=always
