# pkg/production/rell.socket
[Unit]
Description=Rell Socket

[Socket]
ListenStream=0.0.0.0:43800
ListenStream=0.0.0.0:43801

[Install]
WantedBy=sockets.target
